# src/build/modules/registry/files/etc/supervisor/conf.d/registry.conf
[program:registry]
command = /usr/local/bin/docker-registry
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
