# src/registry/build/modules/registry/files/etc/supervisor/conf.d/registry_v2.conf
[program:registry_v2]
command = /go/bin/registry /go/src/github.com/docker/distribution/cmd/registry/config.yml
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
