# src/registry/build/modules/registry/files/etc/supervisor/conf.d/registry.conf
[program:registry]
command = /usr/local/bin/gunicorn --access-logfile - --error-logfile - -k gevent -b 0.0.0.0:5000 -w 4 --max-requests 100 docker_registry.wsgi:application
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
