# src/registry/build/modules/registry/files/etc/profile.d/go.sh
#!/usr/bin/env bash

export GOROOT="/usr/local/src/go"

export PATH="${GOROOT}/bin:${PATH}"

export GOPATH="/go"

export PATH="${PATH}:${GOPATH}/bin"
